# rtcagent.service
[Unit]
Description=Captures packets with eBPF Filter and sends them to Homer
Wants=network.target
After=network.target

[Service]
# Packaged install: see examples/rtcagent.service and /usr/local/rtcagent/etc/rtcagent.env
ExecStart=/usr/local/rtcagent/bin/rtcagent kamailio -m /usr/sbin/kamailio --hep-server="127.0.0.1" --hep-port="9060" --hep-transport="udp" --debug=false
ExecStop=/bin/kill ${MAINPID}
Restart=on-failure
RestartSec=10s
Type=simple

[Install]
WantedBy=multi-user.target
